# power-toggle.service
[Unit]
Description=Apply desktop settings based on the current power source
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=exec
ExecStart=%h/.local/bin/power-toggle monitor
Restart=on-failure
RestartSec=2s
UMask=0077

[Install]
WantedBy=graphical-session.target
